# i3lock.service
[Unit]
Description=I3Lock Lock Manager

[Service]
ExecStart=/usr/bin/i3lock
Restart=always
IgnoreSIGPIPE=no

[Install]
Alias=display-manager.service
